# start-dfs.conf
[program:start-dfs]
command=start-dfs.sh
autostart=true
autorestart=false
exitcodes=0
startretries=0
startsecs=0
priority=2
