# start-yarn.conf
[program:start-yarn]
command=start-yarn.sh
autostart=true
autorestart=false
exitcodes=0
startretries=0
startsecs=0
priority=3
